# docker/app/supervisor/conf.d/api.conf
[program:api]
command=/entrypoint-api.sh
autostart=true
autorestart=false
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0